# init.d/vmcontext
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

name="$SVCNAME"
description="Run OpenNebula contextualization scripts."

: ${mount_dir:="/mnt/context"}
: ${device_label:="CONTEXT"}
: ${script_dirs:="/usr/share/one-context/scripts /etc/one-context.d"}

depend() {
	need localmount
	before netmount
	provide net
}

start_pre() {
	if [ ! -e "$(blkid -l -t LABEL="$device_label" -o device)" ]; then
		eerror "Device with label $device_label does not exist"
		return 1
	fi
}

start() {
	ebegin "Starting $name"

	if ! mountpoint -q "$mount_dir"; then
		mkdir -p "$mount_dir"
		mount -t iso9660 -L "$device_label" -o ro "$mount_dir" || return 1
	fi

	if [ ! -e "$mount_dir"/context.sh ]; then
		eerror "File $mount_dir/context.sh does not exist"
		return 1
	fi

	set -a  # auto-export variables
	. "$mount_dir"/context.sh
	set +a  # disable auto-export

	local dir script
	for dir in $script_dirs; do
		for script in "$dir"/*; do
			[ -x "$script" ] && { "$script" || ewarn "Script $script failed"; }
		done
	done

	umount "$mount_dir"
	rmdir "$mount_dir" 2>/dev/null

	eend 0
}
